# src/services/deepin-service-group@.service.d/dde-file-manage-service-override.conf
# 开放下列安全加固选项给文件管理服务
[Service]
DevicePolicy=
PrivateDevices=no
PrivateNetwork=no
RestrictAddressFamilies=
