# src/services/deepin-service-group@.service.d/dde-file-manage-service-override.conf
# 由于设备管理器对于系统资源的访问，安全加固时需要放开一些权限
[Service]
DevicePolicy=
PrivateDevices=no
PrivateNetwork=no
RestrictAddressFamilies=
